# docker/php/supervisor/tasks/queue.conf
[program:queue]
command=php /var/www/artisan queue:work
numprocs=1
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
